# target/linux/rockchip/armv8/base-files/etc/init.d/uqmi-daemon
#!/bin/sh /etc/rc.common

START=97

start() {
	echo "uqmi daemon is started"
	/usr/bin/uqmi-daemon.sh &
}
